# stow/t430s-trackpoint-sensitivity/etc/udev/rules.d/10-trackpoint.rules
ACTION=="add", KERNEL=="serio2", SUBSYSTEM=="serio", DRIVERS=="psmouse", RUN+="/usr/local/bin/_udev--trackpoint-sensitivity.sh"
