# stow/back-me-up/etc/cron.weekly/rsnapshot_trigger_weekly
#!/bin/sh

TRIGGER=/home/lampam/.back-me-up/WEEKLY

touch "$TRIGGER"
chown lampam:lampam "$TRIGGER"
