# stow/back-me-up/etc/cron.monthly/rsnapshot_trigger_monthly
#!/bin/sh

TRIGGER=/home/lampam/.back-me-up/MONTHLY

touch "$TRIGGER"
chown lampam:lampam "$TRIGGER"
